# config/systemd/.config/systemd/user/udiskie.service
[Unit]
Description=Auto mount removable drives

[Service]
Type=simple
ExecStart=/usr/bin/udiskie -s

[Install]
WantedBy=default.target
