# config/systemd/.config/systemd/user/udiskie.service
[Unit]
Description=Auto mount removable drives

[Service]
Type=Simple
ExecStart=/usr/bin/udiskie -2 -s

[Install]
WantedBy=default.target
